# systemd/.config/systemd/user/battery-notify.service
[Unit]
Description=Send low battery notification

[Service]
Type=oneshot
ExecStart=%h/.config/hypr/scripts/battery-notify.sh

