# systemd/.config/systemd/user/battery-notify.timer
[Unit]
Description=Check battery level every 2 minutes

[Timer]
OnBootSec=1m
OnUnitActiveSec=2m
Unit=battery-notify.service

[Install]
WantedBy=timers.target

